# doc/examples/systemd/ubond@.service.in
[Unit]
Description=UBOND connection to %i
PartOf=ubond.service
ReloadPropagatedFrom=ubond.service

[Service]
Type=notify
NotifyAccess=main
ExecStart=@sbindir@/ubond --config @sysconfdir@/ubond/%i.conf --name %i --user ubond
ExecReload=/bin/kill -HUP $MAINPID
WorkingDirectory=@sysconfdir@/ubond
Restart=on-failure

[Install]
WantedBy=multi-user.target


